# exchange2reminder.service
[Unit]
Description=Exchange to Remind Converter

[Service]
User=<user>
ExecStart=<path to ex2rem.py file> -c <path to ex2rem.conf file>
